# vm-audit.ini
[program:vm-audit]
command=/usr/bin/vm-audit
autostart=true
autorestart=true
